# roles/afilias-os/templates/etc/profile.d/afilias.sh.j2
HISTSIZE={{ shell_history_size }}
HISTTIMEFORMAT="%F %T - "
TMOUT={{ shell_session_timeout }}
PS1="[\u@\h({{ datacenter }}) \W]\\$ "
export HISTSIZE HISTTIMEFORMAT TMOUT PS1
