# roles/cis/files/etc/cron.daily/unowned_files
#!/bin/bash
find / \( \! -fstype xfs -prune \) -o \( \( -nouser -exec chown root {} \; \) , \( -nogroup -exec chgrp root {} \; \) \)
